# Bash-Shell-Script/crontab
MAILTO=\"\"
00 23 * * * root /scripts/backup-aws-S3.sh